# systemd/hedge_fund_ai.service
[Unit]
Description=Hedge Fund AI System Loop
After=network.target

[Service]
WorkingDirectory=/app
ExecStart=/usr/bin/python3 /app/runloop.py
Restart=always
User=ubuntu

[Install]
WantedBy=multi-user.target